# systemd/enrollment.service
[Unit]
Description=Subterra WG coordinator enrollment service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=subterra-hub
Group=subterra-hub
Environment=SUBTERRA_HOST=0.0.0.0
Environment=SUBTERRA_PORT=8443
Environment=SUBTERRA_DB=/var/lib/subterra-hub/state.db
ExecStart=/opt/subterra-hub/.venv/bin/subterra-hub-serve
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/subterra-hub
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
