# systemd/subterra-backup.timer
[Unit]
Description=Nightly backup of the Subterra coordinator

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
